# stage3_UDAS/03-software/files/udas03.service

[Unit]
Description=UDAS03 Service

[Service]
Type=simple
ExecStart=/usr/lib/amplituda/udas03/linux-arm/Amplituda.UDAS.Software.UDAS03.Server
RemainAfterExit=yes
Restart=on-failure

[Install]
WantedBy=multi-user.target